# config/com.sandbox.snapshot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sandbox.snapshot</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>REPO_ROOT_PLACEHOLDER/scripts/snapshot.sh</string>
        <string>sb-dev</string>
    </array>
    <key>StartInterval</key>
    <integer>3600</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>WorkingDirectory</key>
    <string>REPO_ROOT_PLACEHOLDER</string>
    <key>StandardOutPath</key>
    <string>HOME_PLACEHOLDER/.sandbox/snapshots.log</string>
    <key>StandardErrorPath</key>
    <string>HOME_PLACEHOLDER/.sandbox/snapshots.log</string>
    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
